# .config/systemd/user/default.target
[Unit]
Description=Target for user services
