# .config/systemd/user/eclimd.service
[Unit]
Description=Eclim daemon - headless eclipse server

[Service]
ExecStart=/usr/share/eclipse/eclimd

[Install]
WantedBy=default.target
