# .config/systemd/user/urxvtd.service
[Unit]
Description=Urxvt Terminal Daemon

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=default.target
